# cron/sargraph.cron
0 1 * * * root /usr/share/sargraph/cron/parsegraph.sh
