# cron/sargraph.cron
0 1 * * * root /usr/local/sbin/parsegraph.sh
